# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root

[program:llama-server]
command=/opt/llama-venv/bin/python3 -m llama_cpp.server --model /models/Qwen3.5-0.8B-Q4_K_M.gguf --host 127.0.0.1 --port 8080 --n_ctx 2048 --n_threads 4 --verbose false
environment=PORT="8080"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startsecs=10
startretries=3
priority=1

[program:nextjs]
command=node server.js
directory=/app
environment=NODE_ENV="production",HOSTNAME="0.0.0.0",PORT="3000",LLM_BASE_URL="http://127.0.0.1:8080/v1",LLM_MODEL="qwen3.5-0.8b",USE_MOCK_LLM="false",ENCLAVE_MODE="false"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startsecs=3
startretries=3
priority=2
